# data/gdrive-sync-daemon.service
[Unit]
Description=Google Drive sync daemon (gdrive-sync)
After=graphical-session.target

[Service]
Type=dbus
BusName=io.github.teopost.GDriveSync.Daemon
ExecStart=/usr/bin/gdrive-sync-daemon
Restart=on-failure
RestartSec=10
Slice=background.slice

[Install]
WantedBy=default.target
